# aibot.service
[Unit]
Description=AI Neural Networks Telegram Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/home/botuser/ai_neural_networks_bot
Environment=PATH=/home/botuser/ai_neural_networks_bot/venv/bin
ExecStart=/home/botuser/ai_neural_networks_bot/venv/bin/python main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aibot

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/botuser/ai_neural_networks_bot

[Install]
WantedBy=multi-user.target